# api-zvolta.service
[Unit]
Description=Flask HTTP to HTTPS Proxy Server
After=network.target
Documentation=https://github.com/fs0711/http-https-gateway

[Service]
Type=notify
User=root
Group=www-data

WorkingDirectory=/home/gateway
Environment="PATH=/home/gateway/env/bin"
ExecStart=/home/gateway/env/bin/gunicorn --config gunicorn.conf.py wsgi:app

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Process management
KillMode=mixed
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/gateway/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
