# proxy.service
[Unit]
Description=Bidirectional HTTPS/HTTP Proxy
Documentation=file:///opt/proxy/README.md
After=network.target

[Service]
Type=simple
User=proxy
Group=proxy
WorkingDirectory=/opt/proxy

# Environment variables for Let's Encrypt SSL certificates
Environment="FLASK_ENV=production"
Environment="FLASK_DEBUG=False"
Environment="GATEWAY_HOST=0.0.0.0"
Environment="GATEWAY_PORT=5443"
Environment="SSL_ENABLED=True"
Environment="SSL_CERT_PATH=/etc/letsencrypt/live/smartswitch.orkofleet.com/fullchain.pem"
Environment="SSL_KEY_PATH=/etc/letsencrypt/live/smartswitch.orkofleet.com/privkey.pem"
Environment="PROXY_ENDPOINT_A=https://smartswitch.orkofleet.com"
Environment="PROXY_ENDPOINT_B=https://api.zvolta.com"
Environment="PROXY_TIMEOUT=30"
Environment="LOG_LEVEL=INFO"
Environment="LOG_FILE=/var/log/proxy/gateway.log"

# Start command
ExecStart=/opt/proxy/venv/bin/python /opt/proxy/app.py

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/proxy /opt/proxy/logs

# Standard output and error
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxy

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
